# scripts/jerboa/systemd/user/jerboa-market-health-alert.service
[Unit]
Description=Jerboa Market Health - M43 alert run-once service
Documentation=https://github.com/chrisdudley-dev/market-health-cli
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/bin/jerboa-market-health-alert-run-once
TimeoutStartSec=10min

# Logs go to journald by default for user services.
SyslogIdentifier=jerboa-market-health-alert

# Be gentle to the Pi.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

# User-service-compatible hardening.
# Avoid root-only sandboxing directives here; some Raspberry Pi user managers
# fail before ExecStart with status=218/CAPABILITIES.
NoNewPrivileges=yes
PrivateTmp=yes
