# scripts/jerboa/systemd/user/jerboa-market-health-refresh-all-failure.service
[Unit]
Description=Jerboa Market Health - refresh-all failure handler

[Service]
Type=oneshot
ExecStart=%h/bin/jerboa-market-health-alert
