# scripts/jerboa/systemd/user/jerboa-market-health-alert.timer
[Unit]
Description=Jerboa Market Health - run M43 alert service every 15 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m
Persistent=true
Unit=jerboa-market-health-alert.service

[Install]
WantedBy=timers.target
